# main/cron.daily/housekeeping
#!/bin/sh

find /home -mtime +7 -exec rm {} +
